# resources/openrc/alcove
#!/sbin/openrc-run
# Copyright 2020 Bio::Neos, Inc.
# Distributed under the terms of the GNU General Public License v2

extra_commands="checkconfig"

: ${ALCOVE_DIR:=${RC_PREFIX%/}/usr/share/alcove}
: ${ALCOVE_CONFDIR:=${RC_PREFIX%/}/etc/alcove}
: ${ALCOVE_CONFIG:=${ALCOVE_CONFDIR}/alcove.ini}
: ${ALCOVE_PIDFILE:=${RC_PREFIX%/}/run/${SVCNAME}.pid}
pidfile="${ALCOVE_PIDFILE}"

export NODE_ENV=production
# TODO: Don't hardcode log dir! Maybe use a syslogger module?
command="/usr/bin/forever"
command_args="-a -l /var/log/alcove/alcove-forever.log --pidfile ${pidfile} ${ALCOVE_DIR}/app.js"
command_background="true"

depend() {
  use logger net dns
}

checkconfig() {
  if [ ! -e "${ALCOVE_CONFIG}" ] ; then
    eerror "You need an ${ALCOVE_CONFIG} file to run the Alcove Backup System."
    eerror "There is a sample file in ${ALCOVE_CONFDIR}/alcove.ini.example"
    return 1
  fi
}

start_pre() {
  # Sanity check 
  checkconfig || return $?
}

stop_pre() {
  # Sanity check only performed on restarts
  # (i.e. don't stop the daemon if we know it will not be able to restart)
  if [ "${RC_CMD}" = "restart" ] ; then
    checkconfig || return $?
  fi

  # We have to shutdown the alcove script before we end forever:
  ${command} stop ${ALCOVE_DIR}/app.js
}
